# debian/galileo-terremoti.postrm
#!/bin/sh 

set -e

update-rc.d proftpd remove >/dev/null || exit $?

#DEBHELPER#
